# roles/configure_backups/templates/vault_backup.timer.j2
[Unit]
Description=Regularly backup the vault database

[Timer]
OnCalendar={{ bbcrd_vault_backup_schedule }}

[Install]
WantedBy=timers.target
